# ansible/qualifiers/post/graylog.bak/roles/install/templates/opensearch.service.j2
[Service]
Environment="OPENSEARCH_INITIAL_ADMIN_PASSWORD={{ opensearch_initial_admin_password }}"